# conky/.config/autostart/conky.desktop
[Desktop Entry]
Type=Application
Name=Conky
Exec=/home/kugelblitz/.config/conky/start.sh
Hidden=false
NoDisplay=false
X-GNOME-Autostart-enabled=true
